# moodle_cron/templates/minutely-moodle.cron
* * * * * root /usr/bin/php /var/www/html/enrol/database/cli/sync.php
* * * * * root /usr/bin/php /var/www/html/admin/cli/cron.php
