# cookbooks/moodle_opsworks/templates/hourly-moodlebackup.cron.erb
0 * * * * ec2-user /bin/sh /home/ec2-user/backup-moodledata.sh