# cookbooks/moodle_opsworks/templates/moodlebackup.cron.erb
15 0,4,8,12,16,20 * * * ec2-user /bin/bash /home/ec2-user/backup-moodledata.sh
30 0,4,8,12,16,20 * * * ec2-user /bin/bash /home/ec2-user/backup-moodle-mysql.sh
